# rc.local
#!/bin/bash

# wait until i2c is online
while [ ! -e /dev/i2c-1 ]; do
    sleep 1
done

chmod 666 /dev/i2c*
runuser -l ryan -c /home/ryan/room-thermostat/launch_tmux.sh
